# deploy/systemd/cyborg-screen-off.service
[Unit]
Description=Turn off Cyborg kiosk screen overnight

[Service]
Type=oneshot
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/cyborg/.Xauthority
ExecStart=/bin/sh -lc 'xset dpms force off || xrandr --output HDMI-1 --off || tvservice -o'
User=cyborg
Group=cyborg
