# deploy/systemd/cyborg-screen-on.service
[Unit]
Description=Turn on Cyborg kiosk screen in the morning

[Service]
Type=oneshot
Environment=DISPLAY=:0
Environment=XAUTHORITY=@CYBORG_KIOSK_HOME@/.Xauthority
ExecStart=/bin/sh -lc 'xrandr --output HDMI-1 --auto --rotate left || tvservice -p; xset dpms force on'
User=@CYBORG_KIOSK_USER@
Group=@CYBORG_KIOSK_USER@
